# systemd/fandynamic.service
[Unit]
Description=Dell R730XD Dynamic Fan Control Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/root/fandynamic-stable.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fandynamic

[Install]
WantedBy=multi-user.target
